# recipes-data/voyager-data/files/voyager.service
[Unit]
Description=Voyager data startup script

[Service]
ExecStart=/usr/bin/voyager-data.sh

[Install]
WantedBy=default.target
