# recipes-data/voyager-data/files/voyager.timer
[Unit]
Description=Voyager data fetch timer

[Timer]
OnBootSec=10s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
